# ios/layout/Library/LaunchDaemons/dev.r58playz.kdeconnectjb.plist
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.r58playz.kdeconnectjb</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>if [ -x /var/jb/usr/local/bin/kdeconnectd ]; then exec /var/jb/usr/local/bin/kdeconnectd; elif [ -x /usr/local/bin/kdeconnectd ]; then exec /usr/local/bin/kdeconnectd; elif [ -x /.jbroot/usr/local/bin/kdeconnectd ]; then exec /.jbroot/usr/local/bin/kdeconnectd; else exit 127; fi</string>
    </array>

    <key>UserName</key>
    <string>mobile</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>JetsamProperties</key>
    <dict>
        <key>JetsamMemoryLimit</key>
        <integer>20</integer>
        <key>JetsamPriority</key>
        <integer>10</integer>
    </dict>
</dict>
</plist>
